# docker/php-fpm/crontab
* * * * * php /var/www/bin/console twoje:polecenie:symfony
# Pamiętaj o dodaniu pustej linii na końcu pliku
